# crontab
# UTC time, runs at daily login reset
1 16 * * * /app/run.sh
